# systemd/retrynote-worker.service
[Unit]
Description=RetryNote Celery Worker
After=network.target postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=exec
User=retrynote
Group=retrynote
WorkingDirectory=/home/retrynote/app/backend
Environment=PYTHONPATH=/home/retrynote/app/backend
ExecStart=/usr/local/bin/doppler run -- /home/retrynote/app/backend/.venv/bin/celery -A app.workers.celery_app worker --loglevel=info --concurrency=2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
